# fastapi-voice.service
[Unit]
Description=FastAPI Voice Generator Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/FastAPI.bak
Environment=PATH=/usr/bin:/usr/local/bin
ExecStart=/usr/bin/python3 start.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 